# mac-server/com.terminal-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terminal-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>ABSOLUTE_PATH_TO_NODE</string>
        <string>ABSOLUTE_PATH_TO_PROJECT/mac-server/src/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>ABSOLUTE_PATH_TO_PROJECT/mac-server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>ABSOLUTE_PATH_TO_PROJECT/mac-server/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>ABSOLUTE_PATH_TO_PROJECT/mac-server/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>SHELL</key>
        <string>/bin/zsh</string>
        <key>TERM</key>
        <string>xterm-256color</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
